# systemd/inventory-web@.service
[Unit]
Description=Inventory Web Server (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=inventory-%i
Group=inventory-%i
EnvironmentFile=/etc/inventory-system/%i.conf
ExecStart=/opt/inventory-system/venv/bin/inventory-system serve ${INVENTORY_PATH} --port ${WEB_PORT}
Restart=on-failure
RestartSec=5s

# Security
NoNewPrivileges=true
PrivateTmp=true
# Note: Add ReadOnlyPaths manually in the deployed service if needed
# Example: ReadOnlyPaths=/home/user/inventory (systemd doesn't support ${INVENTORY_PATH})

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inventory-web-%i

[Install]
WantedBy=multi-user.target
